# scripts/com.tradingsystem.live.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingsystem.live</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/arielbarack/Programming/PT_Cursor/TradingSystem/run.py</string>
        <string>live</string>
        <string>--force</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/arielbarack/Programming/PT_Cursor/TradingSystem</string>
    
    <key>StandardOutPath</key>
    <string>/Users/arielbarack/Programming/PT_Cursor/TradingSystem/logs/live_trading_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/arielbarack/Programming/PT_Cursor/TradingSystem/logs/live_trading_stderr.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
