# scripts/trading-system.service
[Unit]
Description=Kraken Futures Trading System
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/trading-system
Environment="PATH=/opt/trading-system/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/trading-system/.env
ExecStart=/opt/trading-system/venv/bin/python3 /opt/trading-system/run.py live --force
Restart=always
RestartSec=10
StandardOutput=append:/opt/trading-system/logs/live_trading_stdout.log
StandardError=append:/opt/trading-system/logs/live_trading_stderr.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/trading-system

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
