# scripts/trading-system.service
[Unit]
Description=Trading System Live Trading
After=network.target

[Service]
Type=simple
User=trading
WorkingDirectory=/home/trading/TradingSystem
Environment="PATH=/home/trading/TradingSystem/venv/bin"
ExecStart=/home/trading/TradingSystem/venv/bin/python -m src.entrypoints.prod_live
Restart=always
RestartSec=10
StandardOutput=append:/home/trading/TradingSystem/logs/trading.log
StandardError=append:/home/trading/TradingSystem/logs/trading-error.log

[Install]
WantedBy=multi-user.target
